# server/mood-app.socket
[Unit]
Description=mood-app socket

[Socket]
ListenStream=/run/mood-app/socket

[Install]
WantedBy=sockets.target
